# roles/ansible-elk-kibana/templates/etc/systemd/system/kibana.service.j2
[Unit]
Description=Kibana 4

[Service]
Type=simple
User=kibana
Environment=CONFIG_PATH={{ kibana_dir }}/kibana-{{ kibana_minor_version }}-linux-x64/config/kibana.yml
Environment=NODE_ENV=production
ExecStart={{ kibana_dir }}/kibana-{{ kibana_minor_version }}-linux-x64/bin/kibana

[Install]
WantedBy=multi-user.target
