# roles/ansible-mariadb-galera-cluster/templates/etc/systemd/system/mariadb.service.d/max-open-files.conf.j2

[Service]
LimitNOFILE={{ mariadb_max_open_files }}
