# roles/ansible-mariadb-galera-cluster/templates/etc/systemd/system/mariadb.service.d/oom-score-adjust.conf.j2

[Service]
OOMScoreAdjust={{ mariadb_oom_score_adjust }}
