# roles/ansible-ipset/templates/etc/network/if-pre-up.d/firewall.j2
#! /usr/bin/env bash

{{ ansible_managed|comment }}

ipset restore -! < {{ ipset_config_file }}
iptables-restore {{ ipset_iptables_config_file }}
